# gunicorn.conf
[program:gunicorn]
command=/usr/local/bin/gunicorn gogo:app
directory=/usr/src/app
accesslog=/var/log/gunicorn
